# setup/systemd/breadbrich-auto-deploy.service
[Unit]
Description=Breadbrich Engels auto-deploy (poll origin/main and run safe-deploy when behind)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/breadbrich-backups/auto-deploy.sh
StandardOutput=append:/opt/breadbrich/logs/auto-deploy.log
StandardError=append:/opt/breadbrich/logs/auto-deploy.error.log
